# server/camenc.socket
[Unit]
Description=django camenc socket

[Socket]
ListenStream=/run/camenc.socket

[Install]
WantedBy=sockets.target

